# cloud/registry/transitive-registry.service
[Unit]
Description=Transitive Robotics Registry Service
After=network.target

# Cannot use env variables like $USER but systemd provides some very useful
# specifies: https://www.freedesktop.org/software/systemd/man/systemd.unit.html#Specifiers
[Service]
ExecStart=node index.js
WorkingDirectory=%h/opt/registry
# EnvironmentFile=%h/opt/registry/env
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
