# robot-agent/transitive-package@.service
[Unit]
Description=Transitive Robotics Package Service: %i

[Service]
# System call sandboxing
SystemCallArchitectures=native
NoNewPrivileges=true

WorkingDirectory=%h/.transitive
EnvironmentFile=%h/.transitive/.env
# Env file for user to add other things, this file will not be overwritten
EnvironmentFile=-%h/.transitive/.env_user
Environment=TRPACKAGE=%i

# generate a random password for this package to use
ExecStartPre=%h/.transitive/usr/bin/node -e "fs.writeFileSync('%h/.transitive/packages/%i/password', Math.random().toString(36).substr(2, 9))"

ExecStart=%h/.transitive/unshare.sh "cd /home/transitive && (. /opt/ros/*/setup.bash || true) && PATH=/home/usr/bin:$PATH && . /home/etc/env_local && npm update --no-save && cd node_modules/@transitive-robotics/%i && env PASSWORD=$(cat ../../../password) npm start"

# Note: npm update also installs missing packages, see,
# https://stackoverflow.com/a/19824154/1087119
