# robot-agent/transitive-package@.service
[Unit]
Description=Transitive Robotics Package Update Service

[Service]

# Mount (i.e., filesystem) sandboxing
PrivateUsers=true
ProtectSystem=strict
ProtectHome=tmpfs
BindPaths=%h/.transitive/packages/%i:/home/transitive
BindReadOnlyPaths=%h/.transitive/usr:/home/usr

# System call sandboxing
SystemCallArchitectures=native
NoNewPrivileges=true

WorkingDirectory=/home/transitive/
#ExecStart=/usr/bin/bash -c "/home/usr/bin/node /home/usr/bin/npm update && cd node_modules/%i && /home/usr/bin/node /home/usr/bin/npm start"
ExecStart=/usr/bin/bash -c "PATH=/home/usr/bin:$PATH && npm update && cd node_modules/%i && npm start"
# Note: npm update also installs missing, see,
# https://stackoverflow.com/a/19824154/1087119
